# machines/MALAK/90-yubikey-unplug.rules
ACTION=="remove", ENV{ID_VENDOR_ID}=="1050", ENV{ID_MODEL_ID}=="0407", RUN+="/usr/bin/systemctl start --no-block i3lock.service"
